# systemd/matrixdisplay-daemon.service
[Unit]
Description=matrixdisplay LED matrix daemon
After=local-fs.target

[Service]
Type=simple
WorkingDirectory=/home/pi/matrixdisplay
ExecStart=/home/pi/matrixdisplay/.venv/bin/python -m daemon.run
Restart=on-failure
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
